# provision/service/symmetric.conf
[program:symmetricds]
command= /usr/local/bin/symmetric-server/bin/sym
autorestart=true
stdout_events_enabled=true
stderr_events_enabled=true
